# deploy/systemd/stinger.service
[Unit]
Description=Stinger - Face Recognition with Theme Playback
After=network.target

[Service]
Type=simple
User=stinger
Group=stinger
WorkingDirectory=/opt/stinger/backend
Environment=PATH=/opt/stinger/venv/bin:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
Environment=LD_LIBRARY_PATH=/usr/local/cuda/lib64:/usr/lib/x86_64-linux-gnu
ExecStart=/opt/stinger/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

# Audio access
SupplementaryGroups=audio

# Device access for webcam
SupplementaryGroups=video

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stinger

[Install]
WantedBy=multi-user.target

